# systemd/PrismaConcatDaily.service
[Unit]
Description=Run concatenation of Prisma data

[Service]
ExecStart=/bin/bash PROJECT_PATH/concat.sh

[Install]
WantedBy=multi-user.target